# deploy/cloud-social/cloud-social.service
[Unit]
Description=SlayTheList Cloud Social Service
After=network.target

[Service]
Type=simple
User=slaythelist
WorkingDirectory=/opt/slaythelist
EnvironmentFile=/etc/slaythelist/cloud-social.env
ExecStart=/usr/bin/npm --workspace @slaythelist/cloud-social run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
